# ansible/templates/swag_supervisor.conf.j2
[program:{{app_name}}]
command = /home/{{ansible_ssh_user}}/swag/bin/gunicorn app:app -w4
directory = /home/{{ansible_ssh_user}}/swag/lib/python3.4/site-packages/{{app_name}}
